# claude/.config/systemd/user/claude-notify.service
[Unit]
Description=Claude cross-user notification watcher
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/claude-notify-watch
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
